# target/linux/mediatek/mt7622/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

# SPDX-License-Identifier: GPL-2.0-only

START=99

boot() {
	case $(board_name) in
	belkin,rt3200)
		mtd erase senv
		;;
	esac
}
